# launchagents/com.claude.auto-save-conversations.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.auto-save-conversations</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{CLAUDE_HOME}}/scripts/auto_save_conversations.py</string>
    </array>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Documents/.auto_save.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Documents/.auto_save_error.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
